# templates/freebsd/ftbserver.erb
#!/bin/sh
#
# PROVIDE: <%= @ftb_name %>
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="<%= @ftb_name %>"
rcvar=<%= @ftb_name %>_enable

load_rc_config $name

extra_commands="restart"
start_cmd="${name}_start"
stop_cmd="${name}_stop"
restart_cmd="${name}_restart"

ftb_server_home=<%= @ftb_server_home %>
ftb_server_name=<%= @ftb_server_name%>
ftb_user=<%= @ftb_user %>
ftb_world_name=<%= @ftb_world_name %>

eval "${rcvar}=\${${rcvar}:-'NO'}"

<%= @ftb_name %>_start() {
    if $(/usr/bin/su ${ftb_user} -c 'tmux has-session -t <%= @ftb_name %>'); then
        echo "Server already running."
    else
        echo "Starting ${ftb_server_name}..."
        cd ${ftb_server_home}
        /usr/bin/su ${ftb_user} -c "/usr/local/bin/tmux new-session -d -s <%= @ftb_name %> \"./ServerStart.sh\"" || exit 5

        i=0
        while [ $i -lt 15 ]; do
            i=$(($i + 1))
            if $(/usr/bin/su ${ftb_user} -c 'tmux has-session -t <%= @ftb_name %>'); then
                echo "${ftb_server_name} running."
                break
            else
                sleep 1
            fi
        done

        if [ $i -ge 10 ]; then
            echo "ERROR: Server could not be started."
            /usr/bin/su ${ftb_user} -c "/usr/local/bin/tmux kill-session -t <%= @ftb_name %>"
            exit 10
        fi
    fi
}

<%= @ftb_name %>_stop() {
    if $(/usr/bin/su ${ftb_user} -c 'tmux has-session -t <%= @ftb_name %>'); then
        echo "Halting ${ftb_server_name}."
        /usr/bin/su ${ftb_user} -c "/usr/local/bin/tmux send -t <%= @ftb_name %>:0.0 stop \"Enter\""

        i=0
        while [ $i -lt 15 ]; do
            i=$(($i + 1))
            if $(/usr/bin/su ${ftb_user} -c 'tmux has-session -t <%= @ftb_name %>'); then
                sleep 1
            else
                echo "${ftb_server_name} halted."
                break
            fi
        done

        if [ $i -ge 10 ];then
            echo "ERROR: Server could not be halted."
            exit 20
        fi
    else
        echo "${ftb_server_name} not running."
    fi
}

<%= @ftb_name %>_restart() {
        <%= @ftb_name %>_stop
        sleep 5
        <%= @ftb_name %>_start
}

run_rc_command "$1"
